# ceno-box/ceno-provisioner/roles/bridge/templates/supervisor-ceno.conf.j2
[group:ceno]
programs=bundle-server

[program:bundle-server]
command=/usr/bin/node bundle-server.js
directory={{ bridge_path }}/bundle-server
user={{ bridgeuser }}
environment=CENOLANG="en"
